# crontab
@hourly python3 /app/rssposter.py /app/
